# src/game_controller/utils/game_controller_package/service/robocup_joystick.service
[Unit]
Description=Booster JoyStick Service
After=network.target

[Service]
Type=simple
User=root
ExecStartPre=/bin/sleep 30
ExecStart=/bin/bash /opt/robocup/service/start_joystick.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
